# etc/sv/vnstatd/run
#!/bin/sh

exec 2>&1
exec vnstatd -n -p /run/vnstatd.pid
