# etc/sv/ddclient/run
#!/bin/bash

[ -r conf ] && . ./conf

exec ddclient -foreground ${EXTRA_ARGS}
