# etc/sv/ddclient/run
#!/bin/bash

DAEMON_OPTS=''
[ -r conf ] && . ./conf

exec ddclient -foreground ${DAEMON_OPTS}
