# etc/sv/ddclient/run
#!/bin/bash

[ -f /etc/conf.d/ddclient ] && . /etc/conf.d/ddclient

exec ddclient -foreground ${EXTRA_ARGS}
