# etc/sv/ntpd/run
#!/bin/sh

# can be run as different user
[ -r conf ] && . ./conf
ntpd_user=${NTPD_USER:-ntp}

exec ntpd -g -u "$ntpd_user":"$ntpd_user" -n 2>&1
