# etc/sv/ntpd/run
#!/bin/sh

DAEMON_OPTS=''
[ -r conf ] && . ./conf

exec 2>&1
exec ntpd -g -n ${DAEMON_OPTS}
