# etc/sv/consolekit/run
#!/bin/sh

if ! [ -d /run/ConsoleKit ]; then
	mkdir -p -m0755 /run/ConsoleKit
fi

# wait for services if available
if [ -e /service/cgmanager ]; then
	sv check /service/cgmanager >/dev/null || exit 1
fi
if [ -e /service/dbus ]; then
	sv check /service/dbus >/dev/null || exit 1
fi

exec 2>&1
exec console-kit-daemon --no-daemon
