# etc/sv/haveged/run
#!/bin/sh

exec 2>&1
exec haveged -F -w 1024 -v 1
