# etc/sv/netdata/run
#!/bin/sh

exec chpst -u netdata /usr/sbin/netdata -D -P /var/run/netdata.pid
