# etc/sv/netdata/run
#!/bin/sh

exec chpst -u netdata netdata -D
