# etc/sv/netdata/run
#!/bin/sh

mkdir -p /var/log/netdata
exec chpst -u netdata netdata -D
