# etc/sv/redis/run
#!/bin/sh -e

CONF="/etc/redis.conf"
sed -ie 's/^daemonize yes/daemonize no/' $CONF

mkdir -p /var/lib/redis

exec 2>&1
exec redis-server
