# etc/sv/connmand/run
#!/bin/sh

DAEMON_OPTS=''
[ -r conf ] && . ./conf

exec 2>&1
exec connmand -n ${DAEMON_OPTS}
