# etc/sv/transmission-daemon/run
#!/bin/sh

DAEMON_OPTS=''
[ -r conf ] && . ./conf

exec 2>&1
exec transmission-daemon --foreground ${DAEMON_OPTS}
