# etc/sv/dnsmasq/run
#!/bin/sh
[ -r conf ] && . ./conf
dnsmasq_user="${DNSMASQ_USER:-dnsmasq}"
exec dnsmasq -k --user="$dnsmasq_user"
