# etc/sv/dnsmasq/run
#!/bin/sh

DAEMON_OPTS='--user=dnsmasq'
[ -r conf ] && . ./conf

exec dnsmasq -k ${DAEMON_OPTS}
