# etc/sv/smbd/run
#!/bin/sh
mkdir -p /var/run/samba
exec smbd -F -S
