# etc/sv/cpupower/run
#!/bin/sh
set -e

exec 2>&1
cpupower
exec chpst -b cpupower pause
