# etc/sv/avahi-daemon/run
#!/bin/sh

# if dbus is enabled wait for it.
if [ -e /service/dbus ]; then
	sv check /service/dbus >/dev/null || exit 1
fi

exec 2>&1
exec avahi-daemon
