# etc/sv/klogd/run
#!/bin/sh

DAEMON_OPTS=''
[ -r ./conf ] && . ./conf

exec 2>&1
exec klogd -n $DAEMON_OPTS
