# etc/sv/lightdm/run
#!/bin/sh

sv check dbus >/dev/null || exit 1
mkdir -p /var/run/lightdm
exec lightdm
