# etc/sv/alsa/run
#!/bin/sh
set -e

exec 2>&1
alsactl restore
exec chpst -b alsa pause
