# etc/sv/wicd/run
#!/bin/sh

if [ -e /service/dbus ]; then
  sv check /service/dbus >/dev/null || exit 1
fi

exec wicd --no-daemon
