# etc/sv/nmbd/run
#!/bin/sh
mkdir -p /var/run/samba
exec nmbd -F -S
