# etc/sv/hwclock/finish
#!/bin/bash
[ -r /etc/rc.conf ] && . /etc/rc.conf
[ -z ${HARDWARECLOCK} ] && HARDWARECLOCK="UTC"
hwclock --adjust ${HARDWARECLOCK:+--$(echo $HARDWARECLOCK |tr A-Z a-z)}
