# etc/sv/dkms/finish
#!/bin/sh

[ -r conf ] && . ./conf

exec dkms stop
