# etc/sv/dbus/run
#!/bin/sh

mkdir -p /var/run/dbus

exec 2>&1
exec dbus-daemon --system --nofork
